# profile.d/screen_wacom_orientation.sh
#!/bin/bash
PATH="/opt/screen_wacom_orientation/bin:$PATH"
